# rc.local
#!/bin/bash
#(Multiple lines of comments removed.)
/opt/bootstrap.sh
exit 0
